# supervisord.conf
[supervisord]
nodaemon=true

[program:http_server]
command=python3 -m http.server 8000
directory=/app

[program:app]
command=python3 api_proxy.py
directory=/app
